# webapp.service
[Unit]
Description=WebApp

[Service]
ExecStart=/bin/sh -c /home/clumsyadmin/webapp/start.sh

[Install]
WantedBy=multi-user.target